# backend/docker/crontab
* * * * * php /backend-app/backend/artisan schedule:run
